# ansible/roles/dfe_rdp/files/dfe-vaapi-check.service
[Unit]
Description=DFE VA-API Hardware Encoding Auto-Detection
Documentation=https://gitlab.gnome.org/GNOME/gnome-remote-desktop/-/merge_requests/294
After=graphical.target
Wants=gnome-remote-desktop.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dfe-vaapi-check
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
